# deploy/maestro.service
[Unit]
Description=MAESTRO 2.0
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nocturnesystems
WorkingDirectory=/opt/maestro
EnvironmentFile=-/etc/maestro/runtime.env
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target