# deploy/maestro-update.timer
[Unit]
Description=Run MAESTRO updater every 5 minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target